# builds/install/arch-specific/linux/firebird.service.in
[Unit]
Description=Firebird Database Server
After=syslog.target network.target

[Service]
User=firebird
Group=firebird
ExecStart=@FB_SBINDIR@/fbguard -daemon -forever
StandardError=syslog
RuntimeDirectory=firebird
RuntimeDirectoryMode=0775

[Install]
WantedBy=multi-user.target
